# ansible/roles/we/templates/we.service.j2
[Unit]
Description=auth

[Service]
WorkingDirectory={{root}}
ExecStart={{java}} -jar {{root}}/{{jar}}
SuccessExitStatus=143
StandardOutput=null

[Install]
WantedBy=multi-user.target